# deployment/supervisor/wedding-saas.conf
[program:laravel-horizon]
process_name=%(program_name)s
command=php /var/www/html/artisan horizon
directory=/var/www/html
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/horizon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
stopsignal=SIGTERM

[program:laravel-scheduler]
process_name=%(program_name)s
command=/bin/bash -c "while true; do php /var/www/html/artisan schedule:run; sleep 60; done"
directory=/var/www/html
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
